# .config/systemd/user/disk-monitor.timer
[Unit]
Description=Daily disk usage check

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target
