# .config/systemd/user/disk-monitor.service
[Unit]
Description=Disk usage monitor with desktop alerts

[Service]
Type=oneshot
ExecStart=/home/milo/.scripts/disk-monitor.sh
