# .config/systemd/user/eww-resume.service
[Unit]
Description=Reopen eww windows after resume from suspend/hibernate
After=graphical-session.target

[Service]
Type=simple
ExecStart=/home/milo/.scripts/eww-resume-watch.sh
Restart=always
RestartSec=5

[Install]
WantedBy=graphical-session.target
