# .config/systemd/user/floating-memory.service
[Unit]
Description=Sway per-app floating/tiling memory for tray apps
Documentation=man:sway-ipc(7)
PartOf=sway-session.target
After=sway-session.target

[Service]
Type=simple
ExecStart=%h/.config/sway/scripts/floating-memory.py
Restart=always
RestartSec=2

[Install]
WantedBy=sway-session.target
