# .config/systemd/user/mal-organizer.path
[Unit]
Description=Watch aMule Incoming for completed downloads (MAL local tracking)

[Path]
PathModified=%h/.aMule/Incoming
Unit=mal-organizer.service

[Install]
WantedBy=default.target
