# .config/systemd/user/night-download.service
[Unit]
Description=One-shot night torrent download (WANDANCE S01) via aria2c, no seeding

[Service]
Type=oneshot
ExecStart=%h/.scripts/night-download.sh
