# .config/systemd/user/notif-dismiss-on-focus.service
[Unit]
Description=Dismiss foot desktop notifications when the terminal regains focus
Documentation=man:sway-ipc(7)
PartOf=sway-session.target
After=sway-session.target

[Service]
Type=simple
ExecStart=%h/.config/sway/scripts/notif-dismiss-on-focus.py
Restart=always
RestartSec=2

[Install]
WantedBy=sway-session.target
