# .config/systemd/user/npm-cache-clean.timer
[Unit]
Description=Monthly npm cache cleanup

[Timer]
OnCalendar=monthly
Persistent=true
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target
