# .config/systemd/user/npm-cache-clean.service
[Unit]
Description=Clean npm cache

[Service]
Type=oneshot
ExecStart=/home/milo/.npm-global/bin/npm cache clean --force
