# .config/systemd/user/podman-prune.timer
[Unit]
Description=Weekly Podman cleanup

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target
