# .config/systemd/user/screenshot-clipboard-notify.service
[Unit]
Description=Screenshot clipboard notification daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.config/sway/scripts/screenshot-clipboard-notify.sh
Restart=always
RestartSec=5

[Install]
WantedBy=sway-session.target
