# .config/systemd/user/screenshot-notify.service
[Unit]
Description=Screenshot save notification daemon
After=graphical-session.target

[Service]
Type=simple
ExecStartPre=/usr/bin/rm -f /tmp/screenshot-notify.lock
ExecStart=%h/.config/sway/scripts/screenshot-notify.sh
Restart=always
RestartSec=5

[Install]
WantedBy=sway-session.target
