# .config/systemd/user/trash-purge.service
[Unit]
Description=Purge trash items older than 30 days

[Service]
Type=oneshot
ExecStart=/usr/bin/trash-empty -f 30
