# .config/systemd/user/trash-purge.timer
[Unit]
Description=Daily trash purge (30-day retention)

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
