# .config/systemd/user/wlsunset-warm.service
[Unit]
Description=wlsunset — forced warm gamma (overrides schedule)
PartOf=graphical-session.target
After=graphical-session.target
Conflicts=wlsunset.service

[Service]
Type=simple
ExecStart=%h/.config/sway/scripts/wlsunset-launch.sh --force-warm
Restart=on-failure
RestartSec=5
